# code/stomp/ros_sim/src/rover_ros/rover_launcher/config/rc.local
#!/bin/sh -e

/usr/local/bin/start_modem.py &

su nvidia /home/nvidia/rover/src/rover_ros/rover_launcher/bin/start.sh

exit 0
